# yocto/meta-masterthesis/recipes-apps/can-logger/files/measured.socket
[Socket]
ListenStream=/var/run/measured.sk

[Install]
WantedBy=sockets.target
